# conf/supervisor.ini
# If using a multi-process concurrency (-w 2), a number in the the 2-4 x NUM_CORES range should be good.
# If you are using threads this number can be higher.
[program:pharmacognosy_web]
command=/home/webmaster/.virtualenvs/pharmacognosy/bin/uwsgi conf/uwsgi.ini
numprocs=1
directory=/home/webmaster/apps/pharmacognosy
autostart=true
autorestart=true
startretries=3
user=webmaster

#[program:pharmacognosy_celeryd]
#command=/home/webmaster/.virtualenvs/pharmacognosy/bin/python manage.py celeryd -B
#numprocs=1
#directory=/home/webmaster/apps/pharmacognosy
#autostart=true
#autorestart=true
#startretries=3
#user=webmaster
#stopwaitsecs=600

#[program:pharmacognosy_celerycam]
#command=/home/webmaster/.virtualenvs/pharmacognosy/bin/python manage.py celerycam
#numprocs=1
#directory=/home/webmaster/apps/pharmacognosy
#autostart=true
#autorestart=true
#startretries=3
#user=webmaster
#stopwaitsecs=600

[group:pharmacognosy]
programs=pharmacognosy_web
#pharmacognosy_api,pharmacognosy_celeryd,pharmacognosy_celerycam
